# build/common/root/supervisor-remotion.conf
[program:remotion]
command=node /config/server.js
stdout_logfile=/config/supervisor-remotion.log
stderr_logfile=/config/supervisor-remotion.err.log
autostart=true
autorestart=true
priority=5
